# extras/systemd/user/xreal-mode-watch.service
[Unit]
Description=Xreal Air Mode Watch (auto-activate VR on SBS switch)
After=graphical-session.target

[Service]
Type=simple
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/usr/lib/vr-headset-init/xreal-mode-watch.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
